# vemmd/vemmd.service.in
[Unit]
Description=VE memory mapping daemon
Requires=ve_peermem.service ve-ived.service
After=ve_peermem.service

[Service]
Type=notify
Environment=LOG4C_RCPATH=@sysconfdir@
Restart=always
ExecStartPre=/bin/rm -f @localstatedir@/vemmd.sock ; \
  /bin/sh -c '/bin/ls -t @localstatedir@/core.* 2> /dev/null | while read c; \
    do file $c | \
      awk -F: \'$0 ~ \"from \'\"\'\"\'@sbindir@/vemmd\" \
        { print $1 }\' ; \
    done |  tail -n +4 | xargs rm -f'
ExecStart=@sbindir@/vemmd -s @localstatedir@/vemmd.sock /dev/ve_peermem
ExecStopPost=/bin/rm -f @localstatedir@/vemmd.sock
KillMode=process
LimitCORE=infinity
WorkingDirectory=@localstatedir@

[Install]
WantedBy=multi-user.target
